# deploy/tradingagents.service
[Unit]
Description=TradingAgents autonomous paper trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/trader
EnvironmentFile=%h/trader/.env
ExecStart=%h/trader/.venv/bin/tradingagents run-bot --mode paper --discover
Restart=on-failure
RestartSec=30
UMask=0077
Environment=PYTHONDONTWRITEBYTECODE=1
NoNewPrivileges=true

[Install]
WantedBy=default.target
